# assets/posts/ustanovka-nextcloud/nextcloud-cron.service
[Unit]
Description=Nextcloud cron.php job

[Service]
User=nextcloud
ExecStart=/usr/local/bin/php -f /home/nextcloud/www/cron.php
